# gpu-exporter.service
[Unit]
Description=Kubernetes AMD GPU Metrics Exporter
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Environment=KUBECONFIG=/etc/kubernetes/kubeconfig
Environment=EXPORTER_PORT=9400
Environment=COLLECTION_INTERVAL=300
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=/opt/gpu-exporter
ExecStart=/usr/bin/python3 -u /opt/gpu-exporter/k8s-gpu-exporter.py
Restart=always
RestartSec=30
StartLimitIntervalSec=0
StartLimitBurst=5

[Install]
WantedBy=multi-user.target
